# app/supervisor/supervisord.conf
[program:app]
command=npm run start
user=node
autorestart=true
autostart=true
stderr_events=true
stdout_events=true